# crontab
@hourly cd backup && git push -u origin master
@daily cd scott/reader && . activate && run
