# desk.crontab
@hourly cd ~/Documents/scott-backup    && git pull scott.thomaslevine.com:backup/.git master
@hourly cd ~/Documents/scott/documents && git pull scott.thomaslevine.com:scott/.git/documents master
@hourly cd ~/Documents/scott/listings  && git pull scott.thomaslevine.com:scott/.git/listings master
