# desk.crontab
@hourly cd ~/Documents/scott/logs       && git pull scott.thomaslevine.com:scott-logs/.git master              && git push
@hourly cd ~/Documents/scott/documents  && git pull scott.thomaslevine.com:scott/.git/modules/documents master && git push
@hourly cd ~/Documents/scott/listings-2 && git pull scott.thomaslevine.com:scott/.git/modules/listings master  && git push

# @hourly cd ~/Documents/scott/reader && . ./activate && v2 && cd ../listings-2 && git add . && git commit . -m @hourly\ commit && git push origin master
