# server.crontab
@daily cd scott/reader && . ./activate && run &> ../client/reader-log.txt
@reboot cd scott && ./serve
@hourly cd scott-logs && git add . && git commit . -m automatic\ log\ commit
