# server.crontab
@daily  cd scott/reader && . activate && run
@reboot cd scott && ./serve
