# server.crontab
@daily cd scott/reader && . ./activate && run
@reboot cd scott && ./serve
@hourly cd scott-logs && git add . && git commit . -m automatic\ log\ commit
