# reader/crontab
@daily cd /home/tlevine/scott/reader && . ./activate && run
